# init/ocspcrl.slice
[Unit]
Description="OCSPCRL Server Instances"

[Slice]
